# fleet.service
[Unit]
Description=fleet

[Service]
ExecStart=/usr/bin/fleetd --config=/etc/fleet/fleet.conf
